# build/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info

[program:gunicorn]
command=gunicorn -c ./build/gunicorn.conf.py TransportMap.wsgi:application
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0